# dot_config/systemd/user/cliphist-text.service
[Unit]
Description=Clipboard history manager for text
Documentation=https://github.com/sentriz/cliphist
PartOf=scroll-session.target
After=graphical-session.target
After=wayland-session@scroll.desktop.target

[Service]
Type=simple
ExecStart=/usr/bin/wl-paste --type text --watch cliphist store
Restart=always
RestartSec=1
Slice=background-graphical.slice
PassEnvironment=WAYLAND_DISPLAY
# Prevent rapid restarts if clipboard is being heavily used
StartLimitBurst=10
StartLimitInterval=60s

[Install]
WantedBy=scroll-session.target
