# dot_config/systemd/user/swaybg.service
[Unit]
Description=Wayland wallpaper daemon
Documentation=man:swaybg(1)
PartOf=scroll-session.target
After=graphical-session.target
After=wayland-session@scroll.desktop.target

[Service]
Type=simple
ExecStart=/usr/bin/swaybg -c '#000000'
Restart=always
RestartSec=1
Slice=background-graphical.slice
PassEnvironment=WAYLAND_DISPLAY
# Allow override via environment file if desired
EnvironmentFile=-%h/.config/uwsm/swaybg.env

[Install]
WantedBy=scroll-session.target
