# dot_config/systemd/user/wl-gammarelay.service
[Unit]
Description=Wayland gamma relay daemon
Documentation=https://github.com/jeremija/wl-gammarelay-rs
PartOf=scroll-session.target
After=graphical-session.target
# Ensure Wayland display is available
After=wayland-session@scroll.desktop.target

[Service]
Type=simple
ExecStart=/usr/bin/wl-gammarelay-rs run
Restart=always
RestartSec=1
Slice=background-graphical.slice
# Environment will be inherited from UWSM
PassEnvironment=WAYLAND_DISPLAY

[Install]
WantedBy=scroll-session.target
