# dot_config/containers/litellm.container
[Unit]
Description=LiteLLM Proxy Service
After=network-online.target llm.network.service llm-postgres.service llm-redis.service
Requires=llm.network.service
Wants=llm-postgres.service llm-redis.service

[Container]
Image=ghcr.io/berriai/litellm:main-stable
AutoUpdate=registry
ContainerName=litellm
Network=llm.network
PublishPort=4000:4000
Volume=%h/.config/containers/systemd/_config/litellm.yaml:/app/config.yaml:ro,Z
Environment=LITELLM_MASTER_KEY=sk-llm-master-key-change-this
Environment=DATABASE_URL=postgresql://llm_admin:llm-postgres-password@llm-postgres/litellm
Environment=REDIS_HOST=llm-redis
Environment=REDIS_PORT=6379
Environment=STORE_MODEL_IN_DB=true
Environment=UI_USERNAME=admin
Environment=UI_PASSWORD=changeme
Environment=USE_AIOHTTP_TRANSPORT=true
EnvironmentFile=-%h/.config/containers/llm.env
Exec=--config /app/config.yaml --port 4000 --num_workers 2
HealthCmd="curl -f http://localhost:4000/health/readiness || exit 1"
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target

---
TODO: check is this is the same as the second? different?
---

[Unit]
Description=LiteLLM Proxy Service
After=network-online.target llm.network.service litellm-db.service litellm-redis.service
Requires=llm.network.service
Wants=litellm-db.service litellm-redis.service

[Container]
Image=ghcr.io/berriai/litellm:main-stable
AutoUpdate=registry
ContainerName=litellm
Network=llm.network
PublishPort=4000:4000
Volume=%h/.config/containers/systemd/_config/litellm.yaml:/app/config.yaml:ro,Z
Environment=LITELLM_MASTER_KEY=sk-llm-master-key-change-this
Environment=DATABASE_URL=postgresql://litellm:litellm-password@litellm-db:5432/litellm
Environment=REDIS_HOST=litellm-redis
Environment=REDIS_PORT=6379
Environment=STORE_MODEL_IN_DB=true
Environment=UI_USERNAME=admin
Environment=UI_PASSWORD=changeme
Environment=USE_AIOHTTP_TRANSPORT=true
# Add your API keys here or use EnvironmentFile
Environment=OPENAI_API_KEY=${OPENAI_API_KEY:-sk-placeholder}
Environment=ANTHROPIC_API_KEY=${ANTHROPIC_API_KEY:-sk-placeholder}
Environment=GEMINI_API_KEY=${GEMINI_API_KEY:-placeholder}
EnvironmentFile=-%h/.config/containers/llm.env
Exec=--config /app/config.yaml --port 4000 --num_workers 2
HealthCmd="curl -f http://localhost:4000/health/readiness || exit 1"
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
