# dot_config/containers/llm-postgres.container
[Unit]
Description=Shared PostgreSQL database for LLM services
After=network-online.target llm.network.service
Requires=llm.network.service

[Container]
Image=docker.io/pgvector/pgvector:pg16
AutoUpdate=registry
ContainerName=llm-postgres
Environment=POSTGRES_USER=llm_admin
Environment=POSTGRES_PASSWORD=llm-postgres-password
Environment=POSTGRES_MULTIPLE_DATABASES=litellm,openwebui,mcp
Network=llm.network
Volume=llm-postgres.volume:/var/lib/postgresql/data:Z
HealthCmd="pg_isready -U llm_admin"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=5
# Create multiple databases on first run
Exec=-c shared_preload_libraries=vector -c max_connections=200

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10
# Initialize multiple databases
ExecStartPost=/bin/bash -c 'sleep 5 && podman exec llm-postgres psql -U llm_admin -c "CREATE DATABASE IF NOT EXISTS litellm;" || true'
ExecStartPost=/bin/bash -c 'podman exec llm-postgres psql -U llm_admin -c "CREATE DATABASE IF NOT EXISTS openwebui;" || true'
ExecStartPost=/bin/bash -c 'podman exec llm-postgres psql -U llm_admin -c "CREATE DATABASE IF NOT EXISTS mcp;" || true'

[Install]
WantedBy=scroll-session.target
