# home/dot_config/containers/mcp.container
[Unit]
Description=Model Context Protocol server for tool integrations
After=network-online.target llm.network.service llm-postgres.service
Requires=llm.network.service
Wants=llm-postgres.service

[Container]
Image=ghcr.io/open-webui/mcpo:latest
AutoUpdate=registry
ContainerName=mcp
Network=llm.network
Volume=%h/.config/containers/systemd/_config/mcp-config.json:/app/conf/config.json:ro,Z
Volume=%h/Documents:/data:ro,Z
Environment=DATABASE_URL=postgresql://llm_admin:llm-postgres-password@llm-postgres/mcp
Exec=--config /app/conf/config.json
HealthCmd="curl -fsSL http://127.0.0.1:8000/time/get_current_time -H 'Content-Type: application/json' -d '{\"timezone\": \"America/New_York\"}' || exit 1"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=5
HealthStartPeriod=5s

[Service]
Slice=llm.slice
TimeoutStartSec=600
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
