# home/dot_config/containers/tika.container
[Unit]
Description=Apache Tika content extraction service
After=network-online.target llm.network.service
Requires=llm.network.service

[Container]
Image=docker.io/apache/tika:latest-full
AutoUpdate=registry
ContainerName=tika
Network=llm.network
Environment=TIKA_CHILD_JAVA_OPTS=-JXmx2g
Environment=TIKA_CHILD_MAX_HEAP=2g
HealthCmd="wget -qO- http://127.0.0.1:9998/tika || exit 1"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=5
HealthStartPeriod=5s

[Service]
Slice=llm.slice
TimeoutStartSec=600
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
