# mkosi.extra/usr/lib/systemd/system/var-mnt-nas.automount
[Unit]
Description=Automount NAS share from BE550 USB storage

[Automount]
Where=/var/mnt/nas
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
