# home/dot_config/systemd/user/gtk-settings.service
[Unit]
Description=Apply GTK and GNOME theme settings
PartOf=scroll-session.target
After=graphical-session-pre.target
# Run before any GTK applications start
Before=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c '\
    gsettings set org.gnome.desktop.interface cursor-theme "Bibata-Modern-Classic" && \
    gsettings set org.gnome.desktop.interface color-scheme "prefer-dark" && \
    gsettings set org.gnome.desktop.interface cursor-size 48 && \
    gsettings set org.gnome.desktop.interface text-scaling-factor 1.0 && \
    gsettings set org.gnome.desktop.interface gtk-theme "Catppuccin-Noir-Standard-Green-Dark" && \
    gsettings set org.gnome.desktop.interface.wm.preferences theme "Catppuccin-Noir-Standard-Green-Dark" && \
    gsettings set org.gnome.desktop.wm.preferences button-layout ""'
# Optionally reload settings on reload signal
ExecReload=/bin/sh -c '\
    gsettings set org.gnome.desktop.interface cursor-theme "Bibata-Modern-Classic" && \
    gsettings set org.gnome.desktop.interface color-scheme "prefer-dark" && \
    gsettings set org.gnome.desktop.interface cursor-size 48 && \
    gsettings set org.gnome.desktop.interface text-scaling-factor 1.0 && \
    gsettings set org.gnome.desktop.interface gtk-theme "Catppuccin-Noir-Standard-Green-Dark" && \
    gsettings set org.gnome.desktop.interface.wm.preferences theme "Catppuccin-Noir-Standard-Green-Dark" && \
    gsettings set org.gnome.desktop.wm.preferences button-layout ""'
Slice=session-graphical.slice

[Install]
WantedBy=scroll-session.target
